# installer/linux/kopano-dagent.init.rhel
#!/bin/bash
#
# kopano-dagent Kopano Groupware Core Delivery Agent
#
# chkconfig: - 86 24
# processname: /usr/sbin/kopano-dagent
# config: /etc/kopano/dagent.cfg
# pidfile: /var/run/kopano/dagent.pid

### BEGIN INIT INFO
# Provides: kopano-dagent
# Required-Start: $local_fs $network $remote_fs $syslog
# Required-Stop: $local_fs $network $remote_fs $syslog
# Should-Start: kopano-server
# Should-Stop: kopano-server
# Short-Description: Kopano Groupware Core Delivery Agent
### END INIT INFO

DAGENTCONFIG=/etc/kopano/dagent.cfg
DAGENTPROGRAM=/usr/sbin/kopano-dagent

# Sanity checks.
[ -x $DAGENTPROGRAM ] || exit 0

# the -d option is required to run LMTP mode deamonized
DAGENTCONFIG_OPT="-d"
[ ! -z $DAGENTCONFIG -a -f $DAGENTCONFIG ] && DAGENTCONFIG_OPT="$DAGENTCONFIG_OPT -c $DAGENTCONFIG"

# Switch to C, like at boot time.
unset LANG LC_CTYPE LC_NUMERIC LC_COLLATE LC_MESSAGES LC_ALL
# Update with user settings.
LC_CTYPE=en_US.UTF-8
[ -f /etc/sysconfig/kopano ] && . /etc/sysconfig/kopano
export LANG LC_CTYPE LC_NUMERIC LC_COLLATE LC_MESSAGES LC_ALL

# Source function library.
. /etc/rc.d/init.d/functions

RETVAL=0
dagent=`basename $DAGENTPROGRAM`
lockfile=/var/lock/subsys/$dagent
pidfile=/var/run/kopano/dagent.pid

start() {
	# Start in background, always succeeds
	echo -n $"Starting $dagent: "
	install -dm0775 -o kopano -g kopano /var/run/kopano
	if [ -f /opt/rh/rh-php56/enable ]; then
		. /opt/rh/rh-php56/enable
	fi
	daemon --pidfile="$pidfile" $DAGENTPROGRAM $DAGENTCONFIG_OPT
	RETVAL=$?
	unset LC_ALL LANG
	echo
	[ $RETVAL -eq 0 ] && touch $lockfile

	return $RETVAL
}

stop() {
	echo -n $"Stopping $dagent: "
	killproc -p "$pidfile" $DAGENTPROGRAM
	RETVAL=$?
	echo
	[ $RETVAL -eq 0 ] && rm -f $lockfile $pidfile

	return $RETVAL
}

restart() {
	stop
	start
}

reload() {
	echo -n $"Restarting $dagent: "
	killproc $DAGENTPROGRAM -SIGHUP
	RETVAL=$?
	echo

	return $RETVAL
}

# See how we were called.
case "$1" in
    start)
		start
		;;
    stop)
		stop
		;;
    status)
		status -p "$pidfile" $dagent
		RETVAL=$?
		;;
    restart|force-reload)
		restart
		;;
    condrestart|try-restart)
		if [ -f ${pidfile} ]; then
			stop
			start
		fi
		;;
    reload)
		reload
		;;
    *)
		echo $"Usage: $dagent {start|stop|status|reload|restart|condrestart|force-reload|try-restart}"
		RETVAL=1
		;;
esac

exit $RETVAL
