# installer/linux/kopano-search.service
[Unit]
Description=Kopano Core Search Engine
Documentation=man:kopano-search(8) man:kopano-search.cfg(5)

[Service]
Type=simple
EnvironmentFile=-/etc/sysconfig/kopano
EnvironmentFile=-/etc/default/kopano
ExecStart=/usr/sbin/kopano-search -F
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
