# installer/linux/kopano-spooler.service
[Unit]
Description=Kopano Core Spooler
Documentation=man:kopano-spooler(8) man:kopano-spooler.cfg(5)

[Service]
Type=simple
Environment=LC_ALL=C LANG=C
ExecStart=/usr/sbin/kopano-spooler -F -c /etc/kopano/spooler.cfg

[Install]
WantedBy=multi-user.target
