# installer/linux/kopano-presence.service
[Unit]
Description=Kopano Core Presence Daemon

[Service]
Type=simple
EnvironmentFile=-/etc/sysconfig/kopano
EnvironmentFile=-/etc/default/kopano
ExecStart=/usr/sbin/kopano-presence -F
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
